# templates/systemd/telegram-gateway.service
[Unit]
Description=Telegram Gateway for Hermes
After=network-online.target hermes.service
Wants=network-online.target

[Service]
Type=simple
User=tg-gateway
Group=tg-gateway
EnvironmentFile=/etc/telegram-gateway/gateway.env
WorkingDirectory=/opt/telegram-gateway
ExecStart=/usr/bin/python3 /opt/telegram-gateway/gateway.py
Restart=on-failure
RestartSec=5s

NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=yes
MemoryDenyWriteExecute=yes
UMask=0077

[Install]
WantedBy=multi-user.target
